# roles/backup_estate/templates/v2e-backup-estate.timer.j2
[Unit]
Description=Run v2e-backup-estate.service on a schedule

[Timer]
OnCalendar={{ backup_estate_schedule }}
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
